# kiro-wecom-bridge.service
[Unit]
Description=Kiro WeCom Bridge
After=network.target

[Service]
Type=simple
User=lucky
WorkingDirectory=/mnt/d/code/yami/kiro-wecom-bridge
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=5
Environment=PATH=/home/lucky/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
